# roles/lego/templates/lego.timer.j2
[Unit]
Description=Timer for lego

[Timer]
OnBootSec={{ lego_timer_on_boot }}
OnCalendar={{ lego_timer_on_calendar }}
Unit=lego.service

[Install]
WantedBy=timers.target
